# templates/default/eye_init.lsb.erb
#!/bin/sh
#
# This file was generated by chef for <%= @node[:fqdn] %>
#
### BEGIN INIT INFO
# Provides: <%= @service_name %>
# Required-Start: $network $remote_fs $syslog $local_fs
# Required-Stop: $network $remote_fs $syslog $local_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: starts <%= @service_name %> by eye
### END INIT INFO

EYE=<%= node["eye"]["bin"] %>
CONFIG_FILE=<%= @config_file %>
SERVICE_NAME=<%= @service_name %>
RUNAS=<%= @user || "root" %>

execute() {
  CMD="$EYE $1 $2"
  if [ "$USER" = "$RUNAS" ]; then
    $CMD
  else
    su $RUNAS -c "$CMD"
  fi
}

case "$1" in
  start)
    echo -n "Loading eye configuration for $SERVICE_NAME"
    execute load $CONFIG_FILE
    execute start $SERVICE_NAME
    ;;
  stop)
    execute stop $SERVICE_NAME
    ;;
  restart)
    execute restart $SERVICE_NAME
    ;;
  status)
    execute info $SERVICE_NAME
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac
